# loki.service.j2
[Unit]
Description=Loki service
After=network.target

[Service]
Type=simple
ExecStart=/opt/loki/loki-linux-amd64 -config.file=/opt/loki/loki-local-config.yaml

[Install]
WantedBy=multi-user.target